# packages/kaali-sentinel/deploy/kaali-sentinel.timer
[Unit]
Description=Run Kaali Sentinel host scan every 30 minutes

[Timer]
# Fire quickly after boot — a compromised host often re-establishes its
# backdoor within a minute of reboot (as the pakchoi attacker did).
OnBootSec=30s
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
